# scripts/pkg/service_templates/opensearch/systemd/etc/systemd/system/opensearch-performance-analyzer.service
[Unit]
Description=OpenSearch Performance Analyzer
PartOf=opensearch.service

[Service]
ExecStart=/usr/share/opensearch/bin/performance-analyzer-agent-cli
Restart=on-failure
User=opensearch
Group=opensearch
Environment=OPENSEARCH_HOME=/usr/share/opensearch
WorkingDirectory=/usr/share/opensearch

[Install]
WantedBy=opensearch.service
